# reporting-v2/ops/ai.rudanek.reporting-v2.hourly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>ai.rudanek.reporting-v2.hourly</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/zsh</string>
      <string>/Users/rudolfkonfal/.openclaw/workspace/reporting-v2/scripts/hourly_refresh.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/rudolfkonfal/.openclaw/workspace/reporting-v2</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>AUTO_PUBLISH_PREVIEW</key>
      <string>1</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>StartInterval</key>
    <integer>3600</integer>

    <key>StandardOutPath</key>
    <string>/Users/rudolfkonfal/.openclaw/workspace/reporting-v2/logs/launchd.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/rudolfkonfal/.openclaw/workspace/reporting-v2/logs/launchd.err.log</string>
  </dict>
</plist>
